# contrib/init/aipgd.conf
description "AIPG Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env AIPGD_BIN="/usr/bin/aipgd"
env AIPGD_USER="aipg"
env AIPGD_GROUP="aipg"
env AIPGD_PIDDIR="/var/run/aipgd"
# upstart can't handle variables constructed with other variables
env AIPGD_PIDFILE="/var/run/aipgd/aipgd.pid"
env AIPGD_CONFIGFILE="/etc/aipg/aipg.conf"
env AIPGD_DATADIR="/var/lib/aipgd"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # aipgd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$AIPGD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run aipgd."
        echo "The setting must appear in $AIPGD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following "
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"AIPG Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$AIPGD_PIDDIR"
    chmod 0755 "$AIPGD_PIDDIR"
    chown $AIPGD_USER:$AIPGD_GROUP "$AIPGD_PIDDIR"
    chown $AIPGD_USER:$AIPGD_GROUP "$AIPGD_CONFIGFILE"
    chmod 0660 "$AIPGD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$AIPGD_PIDFILE" \
    --chuid $AIPGD_USER:$AIPGD_GROUP \
    --exec "$AIPGD_BIN" \
    -- \
    -pid="$AIPGD_PIDFILE" \
    -conf="$AIPGD_CONFIGFILE" \
    -datadir="$AIPGD_DATADIR" \
    -disablewallet \
    -daemon

